# Module2/systemd_configurations/currency_exchange.timer
[Unit]
Description=Run currency exchange script every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
